# examples/systemd/minibox.service
[Unit]
Description=MiniBox local proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/minibox
ExecStart=/usr/local/bin/minibox /etc/minibox/source
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/minibox /var/lib/minibox /var/log/minibox
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
